# scripts/podman/vikiclow.container.in
# VikiClow gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into vikiclow's ~/.config/containers/systemd/
# {{VIKICLOW_HOME}} is replaced at install time.

[Unit]
Description=VikiClow gateway (rootless Podman)

[Container]
Image=vikiclow:local
ContainerName=vikiclow
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{VIKICLOW_HOME}}/.vikiclow:/home/node/.vikiclow:Z
EnvironmentFile={{VIKICLOW_HOME}}/.vikiclow/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
